# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:soldi-api]
command=python -m uvicorn main:app --host 0.0.0.0 --port 3001
directory=/app/soldi-api
autostart=true
autorestart=true
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:soldi-website]
command=node --max-old-space-size=512 server.js
directory=/app
autostart=true
autorestart=true
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:discord-bot]
command=node --max-old-space-size=256 discord-bot.js
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:polymarket-bot]
command=node --max-old-space-size=256 polymarket-bot.js
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:twitter-bot]
command=node --max-old-space-size=256 twitter-bot.js
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:social-tracker]
command=node --max-old-space-size=256 social-tracker.js
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
